# systemd/gamepad_launcher.service
[Unit]
Description=Gamepad Launcher (and cleaner) service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStop=-/usr/bin/sudo systemctl restart xow
