# systemd/steamlink.service
[Unit]
Description=Steam Link - Flatpak (per-user)

[Service]
Type=oneshot
ExecStart=/usr/bin/steamlink
ExecStop=-/usr/bin/systemctl --user stop gamepad_launcher.service

[Install]
WantedBy=gamepad_launcher.service
