# systemd/steam_bigpicture.service
[Unit]
Description=Steam Big Picture Mode

[Service]
Type=oneshot
ExecStart=/usr/games/steam -tenfoot
